# sr/initservices/sr_shutdown.service
[Unit]
Description=Shutdown system after timer trigger

[Service]
Type=oneshot
ExecStart=/sbin/shutdown -h now

