# sr/initservices/sr_shutdown.timer
[Unit]
Description=Trigger sr_shutdown.service after 10 seconds

[Timer]
OnActiveSec=10s
AccuracySec=1s

[Install]
WantedBy=timers.target
