# sr/scripts/tmate_ota/tmate-client.service
[Unit]
Description=tmate Client Agent - Reports tmate sessions to server
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/tmate-client
ExecStart=/usr/bin/python3 /opt/tmate-client/agent.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables (optional, can override config.json)
# Environment="TMATE_SERVER_URL=http://your-server:1777"
# Environment="TMATE_API_KEY=your-api-key"
# Environment="TMATE_CLIENT_CONFIG=/etc/tmate-client/config.json"

[Install]
WantedBy=multi-user.target


